# templates/logrotate.timer.j2
{{ ansible_managed | comment }}
[Unit]
Description=Daily User Logrotate for Picodata {{ cluster_name }}
[Timer]
OnCalendar={{ logrotate_shedule }}
Persistent=true
[Install]
WantedBy=timers.target
